# ops/systemd/emem-server.service.example
# emem-server systemd user unit — example configuration.
#
# Install:
#   cp ops/systemd/emem-server.service.example \
#     ~/.config/systemd/user/emem-server.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now emem-server.service
#
# Adapt the Environment= lines below to your operator setup.

[Unit]
Description=emem.dev — Earth memory protocol HTTP/MCP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/emem

# Bind addresses
Environment=EMEM_BIND=127.0.0.1:5051
Environment=EMEM_DATA=/home/ubuntu/emem/var/emem
Environment=EMEM_DEMOS_DIR=/home/ubuntu/emem/var/demos

# TLS (Let's Encrypt via TLS-ALPN-01). Binding :443 as a non-root user
# requires CAP_NET_BIND_SERVICE on the binary or the process. Two
# deployment patterns:
#
#   USER unit (this template — installed under ~/.config/systemd/user/):
#     The kernel does NOT honour AmbientCapabilities= for user-mode
#     systemd because there is no UID transition for the helper to
#     prime. The working pattern is setcap on the binary, applied by
#     scripts/redeploy.sh after every `cargo build --release` (the
#     strip removes the file cap). Leave the line below commented.
#
#   SYSTEM unit (install under /etc/systemd/system/, run as a system
#   service with User=emem): uncomment the AmbientCapabilities line.
#   The system instance can prime the cap before the UID transition.
#
Environment=EMEM_TLS_DOMAINS=emem.dev,www.emem.dev
Environment=EMEM_TLS_BIND=0.0.0.0:443
Environment=EMEM_TLS_CONTACT=mailto:avijeet@vortx.ai
# Environment=EMEM_TLS_STAGING=1
# AmbientCapabilities=CAP_NET_BIND_SERVICE  # system-mode only — see note above

# Logging
Environment=RUST_LOG=info,rustls_acme=info

# Topic-routing inference backend (ort + tokenizers)
Environment=ORT_DYLIB_PATH=/opt/onnxruntime-1.22.0-cuda12/lib/libonnxruntime.so.1.22.0
Environment=LD_LIBRARY_PATH=/opt/onnxruntime-1.22.0-cuda12/lib:/usr/lib/x86_64-linux-gnu

# GPU sidecar (jepa_v2 + Prithvi + Galileo on CUDA via UDS)
Environment=EMEM_SIDECAR_SOCK=%t/emem/jepa_sidecar.sock
Environment=EMEM_SIDECAR_TIMEOUT_MS=10000

# Google Analytics 4 measurement ID. The repo's web/index.html holds a
# `__EMEM_GA_ID__` placeholder; the responder substitutes this value
# at startup. Unset (or comment out) to disable GA entirely — the
# entire GA <script> block is stripped from the served HTML when
# this variable is unset or empty.
Environment=EMEM_GA_MEASUREMENT_ID=

# WRI Global Drivers of Forest Loss v1.2 (Sims et al. 2025) local cache.
# One-time download from Zenodo (record 15366671, ~282 MB COG) lets the
# materializer sample band 1 (driver class 1..7) via the shared
# `cog::sample_pixel_multi` LZW path through a `file://` URL. Unset to
# fall back to the honest `NotImplemented` disclosure that the per-cell
# fetch is blocked by Zenodo's no-Range serving. Bump this when the v1.3
# release is wired (Zenodo record 19485190).
#
# One-time download:
#   mkdir -p $EMEM_DATA/wri_gdm_cache
#   curl -L -o $EMEM_DATA/wri_gdm_cache/drivers_forest_loss_1km_2001_2024_v1_2.tif \
#     https://zenodo.org/records/15366671/files/drivers_forest_loss_1km_2001_2024_v1_2.tif
Environment=EMEM_WRI_GDM_CACHE=/home/ubuntu/emem/var/emem/wri_gdm_cache/drivers_forest_loss_1km_2001_2024_v1_2.tif

LimitNOFILE=65536
ExecStart=/home/ubuntu/emem/target/release/emem-server
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
